# opal-dynamic-policy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:pep-proxy]
command=uvicorn pep_proxy.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/pep-proxy.log
stderr_logfile=/var/log/supervisor/pep-proxy.err
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:bundle-server]
command=uvicorn bundle_server.main:app --host 0.0.0.0 --port 8001
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/bundle-server.log
stderr_logfile=/var/log/supervisor/bundle-server.err
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
